# com.granola-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.granola-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <!-- UPDATE THIS PATH to wherever you put granola_sync.py -->
        <string>/Users/YOUR_USERNAME/granola-sync/granola_sync.py</string>
    </array>
    <key>StartInterval</key>
    <!-- Runs every 300 seconds (5 minutes) -->
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>/tmp/granola-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/granola-sync.log</string>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
